# deployment/vagrant/apiserver_supervisor.conf
[program:node_server]
command=node /usr/ubshared/nodejs/apiserver.js
autostart=true
autorestart=true
stderr_logfile=/var/log/apiserver.log
stdout_logfile=/var/log/apiserver.log
