# deployment/vagrant/node_server_supervisor.conf
[program:node_server]
command=node /usr/js_webgames/nodejs/node_server.js
autostart=true
autorestart=true
stderr_logfile=/var/log/node_server.log
stdout_logfile=/var/log/node_server.log
